# poke-radar.service
[Unit]
Description=Poke-Radar — Card resale opportunity radar
After=network.target

[Service]
Type=simple
User=poke-radar
WorkingDirectory=/opt/poke-radar
ExecStart=/opt/poke-radar/poke-radar-backend
Restart=always
RestartSec=5
EnvironmentFile=/etc/poke-radar/env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=poke-radar

[Install]
WantedBy=multi-user.target
